# roles/transmission/templates/81-vpn-firewall.rules
KERNEL=="tun0", ACTION=="add", RUN+="/usr/local/bin/change_transmission_bind_address.sh add"
